# Library/LaunchAgents/com.jcrawley.biz-scan.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jcrawley.biz-scan</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jcrawley/.claude/schedules/run-report.sh</string>
        <string>/book-scan</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/jcrawley</string>

    <!-- Monday (1) at 7:00 AM local time — after PB reports at 6:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>WakeForNetworkAccess</key>
    <true/>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/jcrawley/.claude/logs/book-scan.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jcrawley/.claude/logs/book-scan-error.log</string>
</dict>
</plist>
