# Library/LaunchAgents/com.jcrawley.aca-gm-report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jcrawley.aca-gm-report</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>[ $(date +%d) -le 7 ] &amp;&amp; /Users/jcrawley/.claude/schedules/run-report.sh /aca-gm-report</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/jcrawley</string>

    <!-- First Wednesday of each month — fires every Wednesday, guard exits if day > 7 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>3</integer>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>3</integer>
    </dict>

    <key>WakeForNetworkAccess</key>
    <true/>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/jcrawley/.claude/logs/aca-gm-report.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jcrawley/.claude/logs/aca-gm-report-error.log</string>
</dict>
</plist>
